# services/cfengine/cf-agent-hourly.service.tmpl
[Unit]
Description=Hourly cf-agent run
Documentation=man:cf-agent
After=multi-user.target

[Service]
Type=oneshot
ExecStart=/usr/sbin/cf-agent
ExecStart=/usr/local/bin/cf-mailresult
